# prime-registry/deploy/prime-registry.service
[Unit]
Description=Numinal — Prime Registry API
After=network.target

[Service]
Type=simple
User=charles
WorkingDirectory=/home/charles/.openclaw/workspace/prime-registry
Environment="PRIME_REGISTRY_SECRET=4539122fcda848adec1646f33fa1ecfcd295725e9ed3925bb350ab7d966514ac"
Environment="REGISTRY_BASE_URL=https://numinals.io"
Environment="CONTRA_WEBHOOK_SECRET=FILL_IN_FROM_CONTRA_DASHBOARD"
Environment="SMTP_USER=FILL_IN_YOUR_EMAIL"
Environment="SMTP_PASS=FILL_IN_YOUR_PASSWORD"
Environment="SOLANA_NETWORK=mainnet-beta"
ExecStart=/home/charles/.local/bin/uvicorn src.api:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
